# build/crontab
0 0 * * * /download.sh >> /var/log/cron.log 2>&1